# systemd/FebusCleanDaily.service
[Unit]
Description=Run deletion of old DAS dirs 

[Service]
ExecStart=/bin/bash PROJECT_PATH/clean.sh

[Install]
WantedBy=multi-user.target